# scripts/systemd/nodeyez-ipaddress.service
# Nodeyez: systemd unit for IP Address panel generation
# /etc/systemd/system/nodeyez-ipaddress.service

[Unit]
Description=IP Addresses
StartLimitIntervalSec=600
StartLimitBurst=2

[Service]
EnvironmentFile=/etc/nodeyez.conf
WorkingDirectory=/home/nodeyez/nodeyez/scripts
ExecStart=/home/nodeyez/.pyenv/nodeyez/bin/python3 /home/nodeyez/nodeyez/scripts/ipaddress.py
User=nodeyez
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
